# root/etc/services.d/ombi/run
#!/usr/bin/with-contenv bash

cd /opt/ombi
exec \
	s6-setuidgid abc ./Ombi
